# chip_boot_blink.service
[Unit]
Description=CHIP Bootup Blink Notification

[Service]
Type=oneshot
ExecStart=/usr/local/bin/boot_blink.sh

[Install]
WantedBy=getty.target

